# assets/braintwin-chroma-backup.timer
[Unit]
Description=Trigger BrainTwin Chroma backup nightly

[Timer]
# 03:30 UTC daily. 30 min after DLM EBS snapshots at 03:00 to avoid
# io contention between the snapshot freeze and the tar.
OnCalendar=*-*-* 03:30:00
# If the box was off at the trigger time, run on next boot.
Persistent=true

[Install]
WantedBy=timers.target
