# assets/braintwin-chroma-backup.service
[Unit]
Description=BrainTwin nightly Chroma tarball to S3
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/braintwin-chroma-backup.sh
